# rayCluster/ray-head.service
[Unit]
Description=Ray Head Node (Apptainer)
After=network.target

[Service]
Type=simple
ExecStart=/home/unhe/factorySim/rayCluster/ray-head.sh
Restart=always
User=root
WorkingDirectory=/home/unhe/factorySim/

[Install]
WantedBy=multi-user.target